# cron/cron.daily
#!/bin/bash
. REPLACELBPHTMLAUTHDIR/_plugindirs.sh
cd $LBPHTMLAUTHDIR
 ./start.sh --original -login